# docs/data/hunts.json
{"hunts": [{"shortName": "LSASS Access from Non System Account", "name": "TH-0001", "attack": [{"tactics": ["Credential Access"], "technique": "Technique/T1003"}]}, {"shortName": "C2 Beaconing via Standard Protocols", "name": "TH-0002", "attack": [{"tactics": ["Command and Control"], "technique": "Technique/T1071"}]}, {"shortName": "Compromise via external media and devices", "name": "TH-0003", "attack": [{"tactics": ["Initial Access", "Lateral Movement"], "technique": "Technique/T1091"}, {"tactics": ["Exfiltration"], "technique": "Technique/T1052"}]}, {"shortName": "Suspicious network traffic over DNS", "name": "TH-0004", "attack": [{"tactics": ["Command and Control"], "technique": "Technique/T1568"}, {"tactics": ["Command and Control"], "technique": "Technique/T1071"}, {"tactics": ["Exfiltration"], "technique": "Technique/T1048"}]}, {"shortName": "Web Shells", "name": "TH-0005", "attack": [{"tactics": ["Persistence"], "technique": "Technique/T1505"}]}, {"shortName": "Autoruns Analysis", "name": "TH-0006", "attack": [{"tactics": ["Persistence", "Execution"], "technique": "Technique/T1543"}, {"tactics": ["Execution", "Persistence", "Privilege Escalation"], "technique": "Technique/T1053"}, {"tactics": ["Persistence", "Privilege Escalation"], "technique": "Technique/T1546"}, {"tactics": ["Persistence", "Privilege Escalation"], "technique": "Technique/T1037"}, {"tactics": ["Persistence", "Privilege Escalation"], "technique": "Technique/T1547"}]}, {"shortName": "File Share Discovery", "name": "TH-0007", "attack": [{"tactics": ["Collection"], "technique": "Technique/T1039"}]}]}